# route/init.cron
0 1 * * * sh /www/wwwroot/mybjut.liboyan.cn/route/reboot.sh
5 1 * * * sh /www/wwwroot/mybjut.liboyan.cn/nginx/spilt.sh
